# server/systemd/paperclipper-tunnel.service.example
[Unit]
Description=Paper Clipper Cloudflare tunnel
After=paperclipper-server.service network-online.target
Wants=network-online.target

[Service]
ExecStart=/home/<user>/.local/bin/cloudflared tunnel --config /path/to/repo/server/cloudflared/config.yml run paperclipper
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
